# systemd/lota-agent.service.d/10-xdg-runtime.conf.example
# SPDX-License-Identifier: MIT
#
# Drop-in override for lota-agent.service.
#
# The agent's container-aware secondary IPC listener
# (src/agent/main_utils.c::setup_container_listener) is created
# under $XDG_RUNTIME_DIR/lota/lota.sock so a Steam pressure-vessel
# container can reach it through a single PRESSURE_VESSEL_FILESYSTEMS_RW
# bind-mount. lota-agent.service runs as root and inherits no
# XDG_RUNTIME_DIR from a login session, so without this drop-in the
# secondary listener is silently skipped and the only socket the
# Wine/Proton hook can connect to is /run/lota/lota.sock, which
# pressure-vessel does not mount into the container.
#
# Activate by copying to /etc/systemd/system/, replacing 1000 with
# the operator UID, and reloading systemd:
#
#   uid=$(id -u)
#   sudo install -D -m 0644 \
#       systemd/lota-agent.service.d/10-xdg-runtime.conf.example \
#       /etc/systemd/system/lota-agent.service.d/10-xdg-runtime.conf
#   sudo sed -i "s/=1000/=${uid}/g" \
#       /etc/systemd/system/lota-agent.service.d/10-xdg-runtime.conf
#   sudo systemctl daemon-reload
#   sudo systemctl restart lota-agent.socket
#
# `lota-steam-setup --install-systemd-dropin` automates the same
# substitution + install pair for the calling user.
#
# Multi-user note: the drop-in is single-valued. Hosts with more
# than one interactive operator should use the
# `container_listener_uid` config knob in /etc/lota/lota.conf
# instead. Each entry yields a /run/user/<uid>/lota/lota.sock
# listener at startup, capped at LOTA_CONFIG_MAX_CONTAINER_LISTENERS.
# `lota-steam-setup --register-uid` automates the conf-file edit and
# the lota-agent.service restart for the caller's $SUDO_UID.

[Service]
Environment=XDG_RUNTIME_DIR=/run/user/1000
