# systemd/lota-agent.service.d/20-allow-dev-kernel.conf.example
# SPDX-License-Identifier: MIT
#
# Drop-in override for lota-agent.service that enables the
# --insecure-allow-dev-kernel CLI flag.
#
# INSECURE. Only install on a development host whose kernel does
# NOT yet meet the production anti-tamper baseline:
#
#   - lockdown=integrity (or confidentiality) on kernel cmdline
#   - module.sig_enforce=1
#   - IMA appraisal policy active at /sys/kernel/security/ima/policy
#   - SELinux relabel applied to /dev/tpm[0-9]+ (lota_tpm_device_t)
#   - fs-verity on /usr/bin/lota-agent
#
# A host that satisfies all five passes the production gate and
# does not need this drop-in. A host that satisfies none of them
# (typical Fedora 44 / Ubuntu 24.04 default install) is refused
# by bpf_loader_verify_kernel_runtime_hardening() and the agent
# never reaches the BPF attach step. This drop-in downgrades the
# three kernel-cmdline checks (lockdown, module-sig, IMA) plus
# the rootfs checks (SELinux label, fs-verity self-binary) from
# hard fail to warn-and-continue so the agent can run the
# enforcement layer for SDK testing, examples/block-demo, and
# tests/integration/test_bpf_gates.sh.
#
# Activate:
#
#   sudo install -D -m 0644 \
#       systemd/lota-agent.service.d/20-allow-dev-kernel.conf.example \
#       /etc/systemd/system/lota-agent.service.d/20-allow-dev-kernel.conf
#   sudo systemctl daemon-reload
#   sudo systemctl restart lota-agent.service
#
# After installation the agent runs with both
# --insecure-allow-dev-kernel and --insecure-allow-mutable-rootfs;
# journalctl prints a WRN-level "INSECURE: ..." line per relaxed
# check at startup. Production deployments REMOVE this drop-in.

[Service]
ExecStart=
ExecStart=/usr/bin/lota-agent --insecure-allow-dev-kernel --insecure-allow-mutable-rootfs
