# systemd/lota-agent.socket
# SPDX-License-Identifier: MIT
#
# LOTA Agent - systemd Socket Activation Unit
#
# Installation:
#   install -m 644 lota-agent.socket /usr/lib/systemd/system/
#   systemctl daemon-reload
#   systemctl enable --now lota-agent.socket
#
# Copyright (C) 2026 Szymon Wilczek

[Unit]
Description=LOTA Agent IPC Socket
Documentation=https://github.com/szymonwilczek/lota

[Socket]
ListenStream=/run/lota/lota.sock
SocketMode=0660
SocketUser=root
SocketGroup=lota

# connection backlog
Backlog=16

# create /run/lota directory
DirectoryMode=0750
RuntimeDirectory=lota

# remove socket file on stop
RemoveOnStop=yes

Service=lota-agent.service

[Install]
WantedBy=sockets.target
