# selinux/lota-agent.service
# SPDX-License-Identifier: MIT
#
# LOTA Agent Systemd Service Unit
#
# Copyright (C) 2026 Szymon Wilczek
#
# Installation:
#   cp lota-agent.service /usr/lib/systemd/system/
#   systemctl daemon-reload
#   systemctl enable --now lota-agent
#

[Unit]
Description=LOTA - Linux Open Trusted Attestation Agent
Documentation=https://github.com/szymonwilczek/lota

# Start after basic system services
After=network.target
After=tpm2-abrmd.service
After=local-fs.target

# TPM resource manager should be available
Wants=tpm2-abrmd.service

# Fail if no TPM available
# ConditionPathExists=/dev/tpmrm0

[Service]
Type=simple
ExecStart=/usr/bin/lota-agent --mode monitor
ExecReload=/bin/kill -HUP $MAINPID

# Restart policy
Restart=on-failure
RestartSec=5s

# Run as root (for TPM and BPF)
# TODO: dedicated user with capabilities
User=root
Group=root

# Filesystem restrictions
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes

# Read-write paths
ReadWritePaths=/var/lib/lota
ReadWritePaths=/var/log/lota
ReadWritePaths=/run/lota

# Read-only paths
ReadOnlyPaths=/etc/lota
ReadOnlyPaths=/usr/lib/lota

# Device access
DeviceAllow=/dev/tpmrm0 rw
DeviceAllow=/dev/tpm0 rw

# Capabilities
CapabilityBoundingSet=CAP_SYS_ADMIN CAP_SYS_RAWIO CAP_IPC_LOCK CAP_BPF CAP_PERFMON CAP_DAC_OVERRIDE CAP_DAC_READ_SEARCH
AmbientCapabilities=CAP_SYS_ADMIN CAP_BPF CAP_PERFMON

# No new privileges after start
NoNewPrivileges=yes

# Limit syscalls
SystemCallFilter=@system-service @io-event @network-io
SystemCallFilter=bpf perf_event_open
SystemCallFilter=~@privileged @resources

# Namespace isolation
PrivateUsers=no
ProtectKernelTunables=yes
ProtectKernelModules=no
ProtectKernelLogs=no
ProtectControlGroups=yes

# Network (attestation)
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK

# Memory protection
MemoryDenyWriteExecute=no

# Seccomp (BPF requires)
SecureBits=noroot-locked

#
# Resource limits
#
LimitMEMLOCK=64M
LimitNOFILE=1024

#
# Logging
#
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lota-agent

[Install]
WantedBy=multi-user.target
Alias=lota.service
